# data/export/upstart_scaling/process_instances.conf.erb
start on starting <%= app %>-<%= name %>

task

script
  SCALE=$(su - <%= user %> -c 'echo $<%= name.upcase.gsub("-", "_") %>_SCALE')
  [ -z $SCALE ] && SCALE=0

  i=1
  while [ $i -le $SCALE ]; do
    start <%= app %>-<%= name %>-instance INSTANCE=$i
    i=$(( $i + 1 ))
  done
end script
