# data/export/upstart_scaling/process_instances.conf.erb
start on starting <%= app %>-<%= process.name %>

task

script
  <%= app.upcase %>_<%= process.name.upcase %>_SCALE=$(su - <%= user %> -c 'echo $<%= app.upcase %>_<%= process.name.upcase %>_SCALE')
  [ -z $<%= app.upcase %>_<%= process.name.upcase %>_SCALE ] && <%= app.upcase %>_<%= process.name.upcase %>_SCALE=0
  export <%= app.upcase %>_<%= process.name.upcase %>_SCALE

  i=1
  while [ $i -le $<%= app.upcase %>_<%= process.name.upcase %>_SCALE ]; do
    start <%= app %>-<%= process.name %>-instance <%= app.upcase %>_<%= process.name.upcase %>_INSTANCE=$i
    i=$(( $i + 1 ))
  done

end script

